# systemd/vsg-telegram-poller.service
[Unit]
Description=VSG Telegram Long-Polling Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vsg
WorkingDirectory=/home/vsg/vsm_agent
ExecStart=/usr/bin/python3 /home/vsg/vsm_agent/vsg_telegram_poller.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vsg-telegram-poller

[Install]
WantedBy=multi-user.target
